# agent_linux/cw-tcollector.service
[Unit]
Description=CloudWatch TCollector
After=network.target

[Service]
Type=simple
ExecStart=PYTHON AGENT

[Install]
WantedBy=multi-user.target

